# systemd/suoha-health.service
[Unit]
Description=suoha health check

[Service]
Type=oneshot
ExecStart=/opt/suoha/healthcheck.sh
